# flash/lib/log2ram/log2ram-reload.service
[Unit]
Description=reload log2ram

[Service]
Type=oneshot
ExecStart=/bin/systemctl reload log2ram
