# flash/lib/log2ram/log2ram-reload.timer
[Timer]
OnActiveSec=6h
OnUnitActiveSec=6h

[Install]
WantedBy=timer.target
