# workbench/wol-libvirt.service
[Unit]
After=network.target

[Service]
ExecStart=/opt/wol-libvirt.sh
Restart=always
User=root

[Install]
WantedBy=multi-user.target
